# docker-bash-completion/skelton/etc/profile.d/bash_completion.sh
#!/bin/bash

# Use bash-completion, if available
[[ $PS1 && -f /usr/local/share/bash-completion/bash_completion ]] && \
    . /usr/local/share/bash-completion/bash_completion
